# hospicare_worker.conf
[program:hospicare-worker]
process_name=hospicare_worker_%(process_num)02d
command=php /var/www/html/hospicare/artisan queue:listen database --queue=default --sleep=2 --tries=2
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile= /var/www/html/hospicare/storage/logs/worker.log
